# www/p5-Tapper-Reports-Web/pkg/tapper_reports_web_fastcgi_live.pl.rc
#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/bin/tapper_reports_web_fastcgi_live.pl"
daemon_flags="-daemon"

. /etc/rc.d/rc.subr

pexp="perl: ${daemon}${daemon_flags:+ ${daemon_flags}}"

rc_cmd $1
