# net/motsognir/pkg/motsognir.rc
#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/sbin/motsognir"

. /etc/rc.d/rc.subr

rc_cmd $1
