# net/mosquitto/pkg/mosquitto.rc
#!/bin/sh
#
# $OpenBSD: cvsyncd.rc,v 1.1 2012/04/10 19:52:44 phessler Exp $

daemon="/usr/local/sbin/mosquitto -d"
daemon_flags="-c /etc/mosquitto/mosquitto.conf"

. /etc/rc.d/rc.subr

rc_cmd $1
