# net/opentracker/pkg/opentracker.rc
#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/bin/opentracker"
daemon_flags="-f ${SYSCONFDIR}/opentracker.conf"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

mkdir -p /var/opentracker

rc_cmd $1
