# multimedia/gerbera/pkg/gerbera.rc
#!/bin/ksh
#
# $OpenBSD: $

daemon="${TRUEPREFIX}/gerbera"
daemon_user="_gerbera"

. /etc/rc.d/rc.subr

rc_cmd $1
