# mail/dspam/pkg/dspam.rc
#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/bin/dspam"
daemon_flags=" --daemon"

. /etc/rc.d/rc.subr

rc_cmd $1
