# www/p5-Tapper-Reports-Receiver/pkg/tapper-reports-receiver-daemon.rc
#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/bin/tapper-reports-receiver-daemon"

. /etc/rc.d/rc.subr

pexp="perl: ${daemon}${daemon_flags:+ ${daemon_flags}}"

rc_cmd $1
