# games/hedgewars/pkg/hedgewars_server.rc
#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/bin/hedgewars-server"
daemon_user="_hedgewars"

. /etc/rc.d/rc.subr

rc_reload=NO
rc_bg=YES

rc_cmd $1
