# games/netrek/server/pkg/netrekd.rc
#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/bin/netrekd"
daemon_user="_netrekd"

. /etc/rc.d/rc.subr

rc_cmd $1
