# net/icbd/pkg/icbd.rc
#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/sbin/icbd"

. /etc/rc.d/rc.subr

pexp="icbd: icbd.*"

rc_cmd $1
