# security/wireguard-go/pkg/wireguard_go.rc
#!/bin/ksh
#
# $OpenBSD: rc.template,v 1.12 2018/01/11 19:30:18 rpe Exp $

daemon="${TRUEPREFIX}/bin/wireguard-go"
daemon_flags=tun

. /etc/rc.d/rc.subr

rc_bg=NO
rc_reload=NO

rc_cmd $1
