# net/3proxy/pkg/3proxy.rc
#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/bin/3proxy"
daemon_flags="${SYSCONFDIR}/3proxy/3proxy.cfg"

. /etc/rc.d/rc.subr

rc_cmd $1
