# sysutils/mcollective/pkg/mcollectived.rc
#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/sbin/mcollectived"

. /etc/rc.d/rc.subr

pexp=".*ruby.* ${daemon}${daemon_flags:+ ${daemon_flags}}"
rc_reload=NO

rc_cmd $1
