# net/tinc/pkg/tincd.rc
#!/bin/sh
#
# $OpenBSD:$

daemon="${TRUEPREFIX}/sbin/tincd"

. /etc/rc.d/rc.subr


rc_cmd $1
