# audio/shairport-sync/pkg/shairport_sync.rc
#!/bin/sh

daemon="/usr/local/bin/shairport-sync"
daemon_user="_shairport"
daemon_flags="-d"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
